# puppet_src/modules/capdmatchingadapter/files/capdmatchingadapter-job-integration.conf
description "Capd Matching Adapter"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 3 20
exec java -Xdebug -agentlib:jdwp=transport=dt_socket,server=y,suspend=n,address=8002 -javaagent:/opt/seeker/seekerHelper.jar -Djsse.enableSNIExtension=false -jar /opt/capdmatchingadapter/matching-service-adapter.jar server opt/capdmatchingadapter/adapterConfig.yml 2>&1 | logger -t CAPDMATCHINGADAPTER -p daemon.info
