# puppet_src/modules/siti_simulator/files/simulator-job.conf
description "Capd SitiAgri Simulator"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 3 20

/usr/bin/forever start /opt/siti-simulator/server.js -m 1
