# puppet_src/modules/siti_dbmapwms/templates/dbmapwms-job.erb
#!/bin/bash
# chkconfig: 2345 20 80

#description "Capd dbmapwms service"

. /etc/init.d/functions

   start() {
    cd <%= @package_location %>
    ./start.sh
   }

   stop() {
    cd <%= @package_location %>
     ./stop.sh
   }

   case "$1" in
    start)
  start
  ;;
    stop)
  stop
  ;;
    restart)
  stop
  start
  ;;
    *)
         echo "Usage: $0 {start|stop|restart}"
   esac

   exit 0
